# systemd-user/system-health-sensors.service
[Unit]
Description=System-Health sensors feeder (oneshot)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sensors_reader.py
