# systemd-user/system-health-telemetry.service
# ~/.config/systemd/user/system-health-telemetry.service
[Unit]
Description=System-Health — Telemetry snapshot (log + CSV, v2.1)
Documentation=man:systemd.service(5)
After=default.target

[Service]
Type=oneshot
SyslogIdentifier=telemetry-v2.1
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 %h/.local/bin/health-telemetry.py

[Install]
WantedBy=default.target
