# systemd-user/system-health-exporter.service
[Unit]
Description=System-Health → Prometheus textfile exporter

[Service]
Type=oneshot
ExecStart=%h/.local/bin/system-health-exporter.sh
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
SystemCallFilter=@system-service
