# systemd-user/system-health-doctor.service
[Unit]
Description=System-Health Doctor (one-shot checks)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/tempctl doctor
